# deploy/stock-alarm.service
[Unit]
Description=Stock Alarm Flask Application
After=network.target

[Service]
Type=simple
User=onev
Group=onev
WorkingDirectory=/home/onev/stock-alarm
Environment="PATH=/home/onev/.local/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/onev/.local/bin/uv run gunicorn --bind 127.0.0.1:5000 --workers 2 "app:create_app()"
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
